# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

echo "Starting OpenVPN"

exec \
	/usr/local/bin/ovpn_run